# .config/systemd/user/pz-initbackup.service
[Unit]
Description=Project Zomboid Initial Backup Trigger
After=network.target

[Service]
Type=simple
ExecStart=%h/scripts/pz_backup/pz-process-monitor.sh
Restart=on-failure